# www/crontab
*/10 * * * /usr/local/inetstat/bin/traffic.sh
*/10 * * * /usr/local/inetstat/bin/traffic.sh -e
